# systemd/aiagent.service
[Unit]
Description=AIAGENT Telegram AI Coding Agent
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/aiagent
EnvironmentFile=/opt/aiagent/.env
ExecStart=/opt/aiagent/.venv/bin/python /opt/aiagent/bot.py
Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
